# tmp/DEBIAN/postinst
#!/bin/bash

echo "Running post install script"
adduser tree --system
sudo systemctl daemon-reload
sudo service tree start
sudo systemctl enable tree